# scripts/claude-jobs/launchd/com.humblehalal.review-triage.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- C4 review triage — daily 08:30 local. Edit REPO_PATH below. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.humblehalal.review-triage</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>REPO_PATH/scripts/claude-jobs/run-job.sh</string>
    <string>review-triage</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
  <key>StandardOutPath</key><string>REPO_PATH/logs/launchd-review-triage.out</string>
  <key>StandardErrorPath</key><string>REPO_PATH/logs/launchd-review-triage.err</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
