# scripts/claude-jobs/launchd/com.humblehalal.cross-check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- C2 data cross-check — weekly, Mondays 09:30 local. Edit REPO_PATH below. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.humblehalal.cross-check</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>REPO_PATH/scripts/claude-jobs/run-job.sh</string>
    <string>data-cross-check</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>9</integer><key>Minute</key><integer>30</integer></dict>
  <key>StandardOutPath</key><string>REPO_PATH/logs/launchd-cross-check.out</string>
  <key>StandardErrorPath</key><string>REPO_PATH/logs/launchd-cross-check.err</string>
  <key>RunAtLoad</key><false/>
</dict>
</plist>
